# supervisord.conf
[supervisord]
nodaemon = true
loglevel = blather

[program:nginx]
command = /usr/sbin/nginx

[program:clusto-apiserver]
command = /usr/local/bin/clusto-apiserver
